# firmware/system/pigpiod.service
[Unit]
Description=GPIO daemon

[Service]
ExecStart=/usr/bin/pigpiod -n 127.0.0.1
Type=forking

[Install]
WantedBy=basic.target
